# .config/autostart/slideshow-refresh.desktop
[Desktop Entry]
Name=SlideshowRefresh
Exec=/home/pi/refresh.sh
Terminal=false
Type=Application
