# .config/autostart/slideshow.desktop
[Desktop Entry]
Name=Slideshow
Exec=/home/pi/slideshow.sh start
Terminal=false
Type=Application
